# root/etc/s6-overlay/s6-rc.d/svc-kapowarr/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Change ownership of /app
echo "chown /app"
lsiown -R abc:abc /app

# Start main script as abc
s6-setuidgid abc python3 /app/Kapowarr.py